# servio.service
[Unit]
Description=Servio - Lightweight Service Manager
Documentation=https://github.com/yourusername/servio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt
EnvironmentFile=-/opt/.env
ExecStart=/opt/servio -addr :8080 -db /var/lib/servio/data.db
Restart=always
RestartSec=5s
StandardOutput=append:/var/log/servio/servio.log
StandardError=append:/var/log/servio/servio.log
SyslogIdentifier=servio

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
